# configs/com.scai.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.scai.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/o005813/SCAI/.venv/bin/python</string>
        <string>/Users/o005813/SCAI/scripts/daily_pipeline.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/o005813/SCAI</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/o005813/SCAI/src</string>
        <key>DYLD_LIBRARY_PATH</key>
        <string>/Users/o005813/SCAI/.local/lib</string>
        <key>PATH</key>
        <string>/Users/o005813/SCAI/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Run at 22:00 UTC (6pm ET) every weekday (Mon-Fri) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>22</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>22</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>22</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>22</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer>22</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/o005813/SCAI/data/paper_trading/logs/daily_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/o005813/SCAI/data/paper_trading/logs/daily_stderr.log</string>

    <!-- Run on login/wake to catch up missed days (idempotent: skips if already ran today) -->
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
